# build/observium/files/etc/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND -f /etc/apache2/httpd.conf
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cron]
autorestart=true
command=crond -f

[program:observium-bootstrap]
command = /opt/observium-bootstrap.sh
startsecs = 0
autorestart = false
startretries = 1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0